# unitA.service
[Unit]
Description = Unit A

[Service]
Type = simple
ExecStart = /usr/local/sbin/unitA

