# unitB.service
[Unit]
Description = Unit B

[Service]
Type = simple
ExecStart = /usr/local/sbin/unitB

